# local.d/ipv6.start
#!/bin/bash
set -m

ip tunnel add he-ipv6 mode sit remote 216.66.80.30 local 188.1.2.3 ttl 255
ip link set he-ipv6 up
ip addr add 2001:470:1001:123::2/64 dev he-ipv6
ip addr add 2001:470:1002:123::1/64 dev he-ipv6
ip addr add 2001:470:1002:123::2 dev he-ipv6
ip route add ::/0 dev he-ipv6
ip -f inet6 addr

